# vendor/github.com/thisisaaronland/go-chicken/vendor/github.com/cooperhewitt/go-ucd/init.d/ucd-server.sh.example
#!/bin/sh -e

### BEGIN INIT INFO
# Provides:          ucd-server.sh
# Required-Start:    $network $remote_fs $named
# Required-Stop:     $network $remote_fs $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts ucd-server.sh
# Description:       Starts ucd-server.sh using start-stop-daemon
### END INIT INFO

# Don't forget to run:
# update-rc.d ucd-server.sh defaults

UCD_USER=YOUR_UCD_SERVER_USER
UCD_DAEMON=YOUR_UCD_SERVER_PATH
UCD_PORT=YOUR_UCD_SERVER_PORT

UCD_ARGS="--port ${UCD_PORT}"

PID_PATH=/var/run/
PID_NAME='ucd-server'
PID_FILE=${PID_PATH}${PID_NAME}.pid

if [ "$(id -u)" != "0" ]; then
    echo "This script must be run as root" 1>&2
    exit 1
fi

case $1 in
    debug)

        echo "Starting ucd in debug mode w/ options: ${UCD_ARGS}"
        exec sudo -u $UCD_USER $UCD_DAEMON $UCD_ARGS
        ;;
    start)

	if [ -e "${PID_FILE}" ]; then
            echo "Error! ${PID_NAME} is currently running!" 1>&2
            exit 1
	fi

        echo "Starting ucd w/ options: ${UCD_ARGS}"
        PID=`exec sudo -u $UCD_USER $UCD_DAEMON $UCD_ARGS  > /dev/null 2>&1 & echo $!`

        if [ -z $PID ]; then
            printf "%s\n" "Fail"
        else
            echo $PID > $PID_FILE
            printf "%s\n" "Ok"
	    echo $PID
        fi

        ;;
    stop)

	printf "%-50s" "Stopping $PID_NAME"
        PID=`cat $PID_FILE`

        if [ -f $PID_FILE ]; then
            kill -HUP $PID
            printf "%s\n" "Ok"
            rm -f $PID_FILE
        else
            printf "%s\n" "pidfile not found"
        fi

        ;;
    restart)
        $0 stop
        sleep 1
        $0 start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|debug}" >&2
        exit 1
        ;;
esac
